# service/backup.service
[Unit]
Description=Daily backup service

[Service]
ExecStart=/usr/bin/bash ${HOME}/.local/bin/backup.sh
